# supervisor-docker-spark/supervisord.conf.template
[supervisord]
nodaemon=true

[program:sshd]
command=service sshd start
autorestart=false

[program:bash]
command=/bin/bash
