# deploy_scripts/data_conf/master_supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=service sshd start
autorestart=false

[program:bootstrap]
command=/etc/bootstrap.sh -bash >/var/log/bootstrap_log.log 2>/var/log/bootstrap_log.log
