# deploy_scripts/data_conf/worker_supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=service sshd start
autorestart=false

